# config/supervisord.conf
[program:solis_ginsong_server]
command=/usr/local/bin/solis_ginsong_server
user=solis
stdout_logfile=/var/log/solis/solis.log
stderr_logfile=/var/log/solis/solis_error.log


