# config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:solis_ginlong_service]
command=/usr/local/bin/solis_ginlong_service
stdout_logfile=/var/log/solis/solis.log
stderr_logfile=/var/log/solis/solis_error.log


